# com.sso-self-healing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sso-self-healing</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Edit this path to match your install location -->
        <string>/Users/YOUR_USERNAME/.config/sso-self-healing/sso-refresh.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/.config/sso-self-healing/logs/sso-refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.config/sso-self-healing/logs/sso-refresh.err.log</string>
</dict>
</plist>
